# servadmin_backend/api/crontab/crontab.test
# test crontab entry
* * * * * /bin/cat /dev/null
